# deploy/django/etc/supervisor/gunicorn.ini
[program:gunicorn]
command=%(ENV_HOME)s/Envs/profits/bin/gunicorn main:app -k aiohttp.worker.GunicornUVLoopWebWorker -n gunicorn_aiohttp_uvloop -w 5 -b 0.0.0.0:8000
directory=%(ENV_HOME)s/Git/profits
user=root
autostart=true
autorestart=true
redirect_stderr=true
